# crontab
0 0 * * * /bin/bash /srv/shiny-server/cleanup.sh
